# configuracion/plantilla_monitorizacion.service.in
[Unit]
Description=Servicio de monitorizaciĂ³n del servidor del proyecto de sistemas informĂ¡ticos.
After=network.target

[Service]
User=root
WorkingDirectory=${ruta_directorio_script_monitorizacion}
ExecStart=${ruta_script_monitorizacion}
Restart=always

[Install]
WantedBy=multi-user.target
